# deploy/com.musicnewsradar.digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.musicnewsradar.digest</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/theoherve/.local/bin/uv</string>
    <string>run</string>
    <string>python</string>
    <string>-m</string>
    <string>scripts.send_digest</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/theoherve/PycharmProjects/music-news-radar</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/theoherve/Library/Logs/music-news-radar/digest.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/theoherve/Library/Logs/music-news-radar/digest.err.log</string>
</dict>
</plist>
