# deploy/com.musicnewsradar.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.musicnewsradar.pipeline</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/theoherve/.local/bin/uv</string>
    <string>run</string>
    <string>python</string>
    <string>-m</string>
    <string>scripts.run_pipeline</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/theoherve/PycharmProjects/music-news-radar</string>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/theoherve/Library/Logs/music-news-radar/pipeline.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/theoherve/Library/Logs/music-news-radar/pipeline.err.log</string>
</dict>
</plist>
